# debian/zgres-deadman.service
[Unit]
Description=Zgres Database Controller

[Service]
User=root
PrivateTmp=true
ProtectHome=read-only
NoNewPrivileges=true
PrivateDevices=true
ExecStart=/usr/bin/zgres-deadman
Restart=always

[Install]
WantedBy=multi-user.target

